# templates/docker-galaxy.service.erb
[Unit]
Description="Docker container for galaxy"
After=docker.service
Requires=docker.service

[Service]
User=vagrant
TimeoutStartSec=0
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/galaxy-bitwf/scripts/start_bitwf.sh
ExecStop=/usr/local/galaxy-bitwf/scripts/stop_bitwf.sh

[Install]
WantedBy=multi-user.target
